# xorg/udev/99-remap-mouse-button.rules
ACTION=="add", SUBSYSTEM=="input", KERNEL=="mouse*" RUN+="/bin/su -l -c '/home/aj/dotfiles/xorg/remap_mouse_button.sh' aj"

